# ticketgal.service
[Unit]
Description=TicketGal Application Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ticketgal
Group=ticketgal
WorkingDirectory=/home/ticketgal/TicketGal
ExecStart=/usr/bin/bash /home/ticketgal/TicketGal/start-prod.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ticketgal
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security settings
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/ticketgal/TicketGal

[Install]
WantedBy=multi-user.target
